# etc/supervisor/conf.d/core.conf
[program:core-deamon]
priority=15
directory=/
command=core-daemon
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=
redirect_stderr=true